# git_ipc/release/rootfs_16c_v300/etc/udev/rules.d/50-firmware.rules
# Load firmware
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/etc/udev/firmware.sh"